# queries/wsl/08-crontab.sql
-- crontab — user + system crontabs
SELECT * FROM crontab ORDER BY path, minute, hour;
