# queries/wsl/10-startup-items.sql
-- startup_items — login/profile-time startup hooks discovered by osquery
SELECT * FROM startup_items ORDER BY name COLLATE NOCASE;
