# queries/windows/09-startup-items.sql
-- startup_items — Run keys, Startup folders, login items
SELECT * FROM startup_items ORDER BY name COLLATE NOCASE;
